# mint/autostart/a-chrome.desktop
[Desktop Entry]
Type=Application
Name=a-chrome
Comment=
Exec=google-chrome
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false

